# etc/profile.d/10functions
# =============================================================================
# shell functions
# =============================================================================

function hascmd () {
  command -v $* > /dev/null 2>&1
}

# vim: ft=zsh
