# etc/profile.d/10mepath.ps1
# =============================================================================
# paths to my local tools
# =============================================================================

$ME_PATH = "$env:ME\local\bin;$env:ME\bin"
