# etc/profile.d/10mepath
# =============================================================================
# paths to my local tools
# =============================================================================

ME_PATH=~/local/bin:~/bin

export ME_ZSH_FUNCTIONS=$ME_LIBEXEC_DIR/zsh/functions
source $ME_ZSH_FUNCTIONS

# vim: ft=zsh
