# etc/profile.d/20term
# =============================================================================
# terminal settings
# =============================================================================

export TERMINFO=~/.terminfo/$ME_HOST
mkdir -p $TERMINFO

# vim: ft=zsh
