# etc/profile.d/30dotnet
# =============================================================================
# .NET
# -----------------------------------------------------------------------------
# https://www.microsoft.com/net/
# =============================================================================

export DOTNET_ROOT=$ME_PLATFORM_TOOL_DIR/dotnet
ME_DOTNET_VERSIONS=(2.1.811 2.2.207 3.1.201 3.1.403)
ME_PATH+=:$ME/.dotnet/tools:$DOTNET_ROOT

export DOTNET_CLI_TELEMETRY_OPTOUT=1

alias dn='dotnet'

# vim: ft=zsh
