# etc/profile.d/30dotnet
# =============================================================================
# .NET
# -----------------------------------------------------------------------------
# https://www.microsoft.com/net/
# =============================================================================

export DOTNET_ROOT=$ME_PLATFORM_TOOL_DIR/opt/dotnet-sdk

ME_PATH+=:~/.dotnet/tools:$DOTNET_ROOT

# vim: ft=zsh
