# etc/profile.d/30dotnet
# =============================================================================
# .NET
# -----------------------------------------------------------------------------
# https://www.microsoft.com/net/
# =============================================================================

export DOTNET_ROOT=$ME_PLATFORM_TOOL_DIR/dotnet
ME_DOTNET_VERSIONS=(2.1.809 2.2.207 3.1.402)
ME_PATH+=:$ME/.dotnet/tools:$DOTNET_ROOT

alias dn='dotnet'

# vim: ft=zsh
