# etc/profile.d/30dotnet
# =============================================================================
# .NET
# -----------------------------------------------------------------------------
# https://www.microsoft.com/net/
# =============================================================================

ME_PATH+=:$ME/.dotnet/tools

# vim: ft=zsh
