# etc/profile.d/30erlang
# =============================================================================
# Erlang
# -----------------------------------------------------------------------------
# http://www.erlang.org/
# =============================================================================

erps () {
  ps -A | grep beam | grep -v grep
}

# vim: ft=zsh
