# etc/profile.d/30fsharp
# =============================================================================
# F#
# =============================================================================

alias fs='dotnet fsi'

# vim: ft=zsh
