# etc/profile.d/30go
# =============================================================================
# Go
# -----------------------------------------------------------------------------
# https://golang.org/
# =============================================================================

ME_GO_VERSIONS=(1.20 1.19.5)
export GOROOT=$ME_PLATFORM_TOOL_DIR/go/goroot
export GOPATH=$ME_PLATFORM_TOOL_DIR/go/gopath

ME_PATH=$GOROOT/bin:$ME_PATH

# vim: ft=zsh
