# etc/profile.d/30go.ps1
# =============================================================================
# Go
# -----------------------------------------------------------------------------
# https://golang.org/
# =============================================================================

$Env:GOPATH = "$Env:ME_PLATFORM_TOOL_DIR\golang"
