# etc/profile.d/30go
# =============================================================================
# Go
# -----------------------------------------------------------------------------
# https://golang.org/
# =============================================================================

ME_GO_VERSIONS=(1.12.1)
export GOPATH=$ME_PLATFORM_TOOL_DIR/go

ME_PATH+=:$ME_PLATFORM_TOOL_DIR/go/bin

# vim: ft=zsh
