# etc/profile.d/30go
# =============================================================================
# Go
# -----------------------------------------------------------------------------
# https://golang.org/
# =============================================================================

ME_GO_VERSIONS=(1.13 1.12.9)
export GOROOT=$ME_PLATFORM_TOOL_DIR/homebrew/opt/go/libexec
export GOPATH=$ME_PLATFORM_TOOL_DIR/golang

ME_PATH+=:$GOPATH/bin

# vim: ft=zsh
