# etc/profile.d/30go
# =============================================================================
# Go
# -----------------------------------------------------------------------------
# https://golang.org/
# =============================================================================

export GOPATH=$ME_DISTRO_TOOL_DIR/go

# vim: ft=zsh
