# etc/profile.d/30nodejs
# =============================================================================
# Node.js
# =============================================================================

export NPM_HOME=$ME/.npm
ME_PATH+=:$NPM_HOME/packages/bin

# vim: ft=zsh
