# etc/profile.d/30python.ps1
$Env:PYTHONPATH="$Env:ME_LIB_DIR\python"
$Env:PYTHONSTARTUP="$Env:ME_LIBEXEC_DIR\pythonstartup.py"