# etc/profile.d/30ruby
# =============================================================================
# Ruby
# =============================================================================

export GEM_HOME=$ME_DISTRO_TOOL_DIR/lib/ruby/gems
ME_PATH+=:$GEM_HOME/bin

alias rb='ruby'

# vim: ft=zsh
