# etc/profile.d/40pipenv.ps1
$Env:PIPENV_PIPFILE = "$Env:ME\Pipfile"
