# etc/profile.d/40pipenv
# =============================================================================
# Pipenv
# -----------------------------------------------------------------------------
# https://docs.pipenv.org/
# =============================================================================

ME_PIPENV_PIPFILE=$ME/Pipfile

# vim: ft=zsh
