# etc/profile.d/50ant
# =============================================================================
# Apache Ant
# =============================================================================

ME_ANT_VERSIONS=(1.9.9 1.10.1 1.8.4)

export ANT_HOME=$ME_TOOL_DIR/opt/apache-ant

ME_PATH+=:$ANT_HOME/bin

# vim: ft=zsh
