# etc/profile.d/50ant
# =============================================================================
# Apache Ant
# =============================================================================

ME_ANT_VERSIONS=(1.9.6 1.8.4)

ME_PATH+=:$ME_TOOL_DIR/opt/ant/bin

# vim: ft=zsh
