# etc/profile.d/50ant
# =============================================================================
# Apache Ant
# =============================================================================

ME_ANT_VERSIONS=(1.10.7 1.9.14)

export ANT_HOME=$ME_TOOL_DIR/opt/ant

ME_PATH+=:$ANT_HOME/bin

# vim: ft=zsh
