# etc/profile.d/50asciidoctor
# =============================================================================
# AsciiDoctor
# -----------------------------------------------------------------------------
# http://asciidoctor.org/
# =============================================================================

alias adoc='asciidoctor'

# vim: ft=zsh
