# etc/profile.d/50bat.ps1
if (Get-Command "bat" -ErrorAction SilentlyContinue)
{
    Function cat($path) {
        bat -P "$path"
    }
}
