# etc/profile.d/50bat
# =============================================================================
# bat
# -----------------------------------------------------------------------------
# https://github.com/sharkdp/bat
# =============================================================================

if hascmd bat; then
  alias bat='bat --theme=OneHalfDark'
  alias cat='bat'
fi

# vim: ft=zsh
