# etc/profile.d/50bat
# =============================================================================
# bat
# -----------------------------------------------------------------------------
# https://github.com/sharkdp/bat
# =============================================================================

if hascmd bat; then
  alias cat='bat -P'
fi

# vim: ft=zsh
