# etc/profile.d/50bash.ps1
Set-Alias bash $Env:SystemRoot\System32\bash.exe
