# etc/profile.d/50bat.ps1
if (Get-command "bat" -ErrorAction SilentlyContinue)
{
  Set-Alias cat bat
}
