# etc/profile.d/50cherokee
# =============================================================================
# Cherokee web server
# -----------------------------------------------------------------------------
# http://cherokee-project.com/
# =============================================================================

alias cal='cherokee-admin-launcher'

# vim: ft=zsh
