# etc/profile.d/50cli4
# =============================================================================
# HTTPie: command-line HTTP client
# -----------------------------------------------------------------------------
# http://httpie.io/
# =============================================================================

alias cli4="PIPENV_PIPFILE=$ME_PIPENV_PIPFILE pipenv run cli4"

# vim: ft=zsh
