# etc/profile.d/50diff
# =============================================================================
# diff utils
# =============================================================================

# -----------------------------------------------------------------------------
# http://www.colordiff.org/
# -----------------------------------------------------------------------------
if which colordir > /dev/null 2>&1 ; then
    alias diff='colordiff'
fi

# vim: ft=zsh
