# etc/profile.d/50fluxbox
# =============================================================================
# Fluxbox windowmanager
# -----------------------------------------------------------------------------
# http://www.fluxbox.org/
# =============================================================================

export FLUXBOX_LOG=$ME_LOG_DIR/fluxbox.$ME_HOST.log

# vim: ft=zsh
