# etc/profile.d/50gettext
# =============================================================================
# GNU getttext
# -----------------------------------------------------------------------------
# https://www.gnu.org/software/gettext/
# =============================================================================

export LANG=en_US.UTF-8

# vim: ft=zsh
