# etc/profile.d/50gemfire
# =============================================================================
# GemFire Cache Server
# -----------------------------------------------------------------------------
# https://pivotal.io/pivotal-gemfire
# =============================================================================

ME_GEMFIRE_VERSIONS=(9.7.0)
ME_PATH+=:$ME_OPT_DIR/pivotal-gemfire-$ME_GEMFIRE_VERSIONS[1]/bin

# vim: ft=zsh
