# etc/profile.d/50github
# =============================================================================
# GitHub
# =============================================================================

HUB_SOURCES=$ME_SRC_DIR/github.com

alias -g hub=$HUB_SOURCES

# vim: ft=zsh
