# etc/profile.d/50gitlab.ps1
$Env:LAB_SOURCES = "$Env:ME_SRC_DIR\gitlab.eng.vmware.com"

function lab { cd "$Env:LAB_SOURCES" }
