# etc/profile.d/50gitlab
# =============================================================================
# VMware GitLab
# =============================================================================

alias lab='cd ~/src/gitlab.eng.vmware.com/'

# vim: ft=zsh
