# etc/profile.d/50github
# =============================================================================
# GitHub
# =============================================================================

alias hub='cd ~/src/github.com'

# vim: ft=zsh
