# etc/profile.d/50github.ps1
function hub { cd $env:ME_SRC_DIR/github.com }
