# etc/profile.d/50gitlab
# =============================================================================
# VMware GitLab
# =============================================================================

LAB_SOURCES=$ME_SRC_DIR/gitlab.eng.vmware.com

alias lab='cd $LAB_SOURCES'

# vim: ft=zsh
