# etc/profile.d/50gnutls
# =============================================================================
# GnuTLS
# ----------------------------------------------------------------------------
# http://www.gnutls.org/
# =============================================================================

if command -v gnutls-certtool >/dev/null ; then
  alias certtool='gnutls-certtool'
fi

# vim: ft=zsh
