# etc/profile.d/50gnupg
# =============================================================================
# GNU Privacy Guard
# -----------------------------------------------------------------------------
# https://www.gnupg.org/
# =============================================================================

if which gpg2 > /dev/null 2>&1 ; then
  alias gpg=gpg2
fi

# vim: ft=zsh
