# etc/profile.d/50gnutls
# =============================================================================
#
# GnuTLS
# ----------------------------------------------------------------------------
# http://www.gnutls.org/
# =============================================================================

if hascmd gnutls-certtool; then
  alias certtool='gnutls-certtool'
fi

# vim: ft=zsh
