# etc/profile.d/50gnutls
# =============================================================================
# GnuTLS
# ----------------------------------------------------------------------------
# http://www.gnutls.org/
# =============================================================================

type gnutls-certtool >/dev/null 2>&1 && alias certtool='gnutls-certtool'

# vim: ft=zsh
