# etc/profile.d/50gnupg
# =============================================================================
# GNU Privacy Guard
# -----------------------------------------------------------------------------
# https://www.gnupg.org/
# =============================================================================

if hascmd gpg2; then
  alias gpg=gpg2
fi

# vim: ft=zsh
