# etc/profile.d/50google-cloud
# =============================================================================
# Google Cloud SDK
# =============================================================================

ME_GCLOUD_SDK_VERSIONS=(228.0.0)
ME_GCLOUD_SDK_DIR=$ME_PLATFORM_TOOL_DIR/google-cloud-sdk

ME_PATH+=:$ME_GCLOUD_SDK_DIR/bin
if [ -f $ME_GCLOUD_SDK_DIR/completion.zsh.inc ]; then
  . $ME_GCLOUD_SDK_DIR/completion.zsh.inc ]
fi

# vim: ft=zsh
