# etc/profile.d/50grip
# =============================================================================
# Grip : Preview GitHub Markdown
# -----------------------------------------------------------------------------
# https://github.com/joeyespo/grip
# =============================================================================

token=$(token github/grip)
if [[ -n $token ]]; then
  alias grip="grip --user=ccheetham --pass=$token"
fi

# vim: ft=zsh
