# etc/profile.d/50gren
# =============================================================================
# Gren
# -----------------------------------------------------------------------------
# https://github.com/github-tools/github-release-notes
# =============================================================================

GREN_GITHUB_TOKEN=$(github_token gren)

# vim: ft=zsh
