# etc/profile.d/50grip
# =============================================================================
# Grip : Preview GitHub Markdown
# -----------------------------------------------------------------------------
# https://github.com/joeyespo/grip
# =============================================================================

if [[ -n $GITHUB_TOKEN_GRIP ]]
then
    alias grip="grip --user=ccheetham --pass=$GITHUB_TOKEN_GRIP"
fi

# vim: ft=zsh
