# etc/profile.d/50groovy
# =============================================================================
# Groovy
# -----------------------------------------------------------------------------
# http://www.groovy-lang.org/
# =============================================================================

ME_GROOVY_VERSIONS=(2.4.15)
ME_PATH+=:$ME_TOOL_DIR/opt/groovy/bin

# vim: ft=zsh
