# etc/profile.d/50http
# =============================================================================
# HTTPie: command-line HTTP client
# -----------------------------------------------------------------------------
# http://httpie.io/
# =============================================================================

alias http="PIPENV_PIPFILE=$ME_PIPENV_PIPFILE pipenv run http --follow"
