# etc/profile.d/50http
# =============================================================================
# HTTPie: command-line HTTP client
# -----------------------------------------------------------------------------
# http://httpie.io/
# =============================================================================

http() {
  cd $ME && pipenv run http --follow $@
}
