# etc/profile.d/50http.ps1
function Http
{
  pipenv run http $Args
}
