# etc/profile.d/50http
# =============================================================================
# HTTPie: command-line HTTP client
# -----------------------------------------------------------------------------
# http://httpie.io/
# =============================================================================

alias http='http --follow'
