# etc/profile.d/50jetbrains.ps1
Set-Alias jbg goland
Set-Alias jbj idea
Set-Alias jbp pycharm
Set-Alias jbr rider
