# etc/profile.d/50kubernetes.ps1
Function k
{
    Invoke-Expression "kubectl $($Args -join " ")"
}

Function kk
{
    k config "$Args"
}
