# etc/profile.d/50kubernetes
# =============================================================================
# Kubernetes
# =============================================================================

if hascmd microk8s.kubectl >/dev/null; then
  alias kubectl='microk8s.kubectl'
fi

alias kk='k config'

k() {
  local cmd="kubectl $*"
  echo "\$ "$cmd
  eval $cmd
}

# vim: ft=zsh
