# etc/profile.d/50kubernetes.ps1
Function k
{
    & kubectl $Args
}

Function kk
{
    & k config $Args
}
