# etc/profile.d/50less
# =============================================================================
# less
# -----------------------------------------------------------------------------
# http://www.greenwoodsoftware.com/less/
# =============================================================================

alias less="less -RSM~gIsw"

# vim: ft=zsh
