# etc/profile.d/50maven
# =============================================================================
# Apache Maven
# =============================================================================

ME_MAVEN_VERSIONS=(3.6.3)

ME_PATH+=:$ME_TOOL_DIR/opt/maven/bin

function m() {
  cmd="./mvnw $@"
  echo "% $cmd"
  eval $cmd
}

# vim: ft=zsh
