# etc/profile.d/50neovim.ps1
if (Get-command "nvim" -ErrorAction SilentlyContinue)
{
  $Env:EDITOR = "nvim"
  Set-Alias vi $Env:EDITOR
}
