# etc/profile.d/50pipenv
# =============================================================================
# Pipenv
# -----------------------------------------------------------------------------
# https://docs.pipenv.org/
# =============================================================================

export PIPENV_PIPFILE=$ME/Pipfile

# vim: ft=zsh
