# etc/profile.d/50poshgit.ps1
Import-Module "$env:ME_REPO_DIR\posh-git\src\posh-git.psd1"
