# etc/profile.d/50powershell.ps1
Set-PSREadLineOption -HistorySavePath "${Env:XDG_CACHE_HOME}\powershell\history"