# etc/profile.d/50redis
# =============================================================================
# Redis: key-value cache/store
# -----------------------------------------------------------------------------
# http://redis.io/
# =============================================================================

alias redis='redis-server'

# vim: ft=zsh
