# etc/profile.d/50sed
# =============================================================================
# sed util
# =============================================================================

if checkexe gsed > /dev/null ; then
    alias sed='gsed'
fi

# vim: ft=zsh
