# etc/profile.d/50silversearcher
# =============================================================================
# The Silver Searcher
# -----------------------------------------------------------------------------
# https://github.com/ggreer/the_silver_searcher
# =============================================================================

AG_ARGS='-s'      # match case-sensitively
alias ag='ag $AG_ARGS'

# vim: ft=zsh
