# etc/profile.d/50skunkworks
# =============================================================================
# Skunk Works
# =============================================================================

alias skunk='cd $ME_SRC_DIR/skunkworks'

# vim: ft=zsh
