# etc/profile.d/50snap
# =============================================================================
# snap
# =============================================================================

if test -d /snap; then
  ME_PATH+=:/snap/bin
fi

# vim: ft=zsh
