# etc/profile.d/50syft
# =============================================================================
# Syft
# -----------------------------------------------------------------------------
# https://github.com/anchore/syft
# =============================================================================

ME_SYFT_VERSION=0.46.1

# vim: ft=zsh
