# etc/profile.d/50tda
# =============================================================================
# Thread-Dump Analzer
# -----------------------------------------------------------------------------
# https://github.com/irockel/tda
# =============================================================================

ME_TDA_VERSION=2.3.3

# vim: ft=zsh
