# etc/profile.d/50tomcat
# =============================================================================
# Apache Tomcat
# =============================================================================

export ME_TOMCAT_VERSION_8=8.0.41
export ME_TOMCAT_VERSION_7=7.0.75
export ME_TOMCAT_VERSION_6=6.0.48

# vim: ft=zsh
