# etc/profile.d/50tomcat
# =============================================================================
# Apache Tomcat
# =============================================================================

export ME_TOMCAT_VERSION_9=9.0.17
export ME_TOMCAT_VERSION_8=8.5.39
export ME_TOMCAT_VERSION_7=7.0.93

# vim: ft=zsh
