# etc/profile.d/50tomcat
# =============================================================================
# Apache Tomcat
# =============================================================================

export ME_TOMCAT_VERSION_9=9.0.30
export ME_TOMCAT_VERSION_8=8.5.50
export ME_TOMCAT_VERSION_7=7.0.99

# vim: ft=zsh
