# etc/profile.d/50tomcat
# =============================================================================
# Apache Tomcat
# =============================================================================

export ME_TOMCAT_VERSION_8=8.0.33
export ME_TOMCAT_VERSION_7=7.0.69
export ME_TOMCAT_VERSION_6=6.0.45

# vim: ft=zsh
