# etc/profile.d/50visualstudiocode
# =============================================================================
# Visual Studio Code
# -----------------------------------------------------------------------------
# https://code.visualstudio.com/
# =============================================================================

if [[ $ME_OS == osx ]]
then
    alias code=/Applications/Visual\\\ Studio\\\ Code.app/Contents/Resources/app/bin/code
fi


# vim: ft=zsh
