# etc/profile.d/50vim
# =============================================================================
# vim
# -----------------------------------------------------------------------------
# https://www.vim.org/
# =============================================================================

if ! hascmd nvim; then
  if hascmd vim; then
    alias vi=vim
  fi
fi

# vim: ft=zsh
