# etc/profile.d/50vim
# =============================================================================
# vim
# =============================================================================

alias vi=vim

# vim: ft=zsh
