# etc/profile.d/50vim
# =============================================================================
# vim
# -----------------------------------------------------------------------------
# https://www.vim.org/
# =============================================================================

if ! command -v nvim > /dev/null ; then
  if command -v vim > /dev/null ; then
    alias vi=vim
  fi
fi

# vim: ft=zsh
