# etc/profile.d/50vtop
# =============================================================================
# vtop
# -----------------------------------------------------------------------------
# https://github.com/MrRio/vtop
# =============================================================================

alias vtop='vtop --theme gruvbox'

# vim: ft=zsh
