# etc/profile.d/50zfs
# =============================================================================
# ZFS
# =============================================================================

alias zfs='sudo zfs'

# vim: ft=zsh
