# etc/profile.d/50zip
# =============================================================================
# zip
# =============================================================================

# unzip-and-delete
uad() {
  local target=$1
  unzip $target && rm $target
}

# vim: ft=zsh
