# etc/profile.d/60editor
# =============================================================================
# editor settings
# =============================================================================

export EDITOR="vi"
export VISUAL="vi -f"

# vim: ft=zsh
