# etc/profile.d/60xml
# =============================================================================
# XML tooling
# =============================================================================

if [[ -f /usr/local/etc/xml/catalog ]] ; then
    export XML_CATALOG_FILES=/usr/local/etc/xml/catalog
fi

# vim: ft=zsh
