# etc/profile.d/90app-accelerator
# =============================================================================
# App Accelerator
# -----------------------------------------------------------------------------
# Samples: https://github.com/sample-accelerators/
# =============================================================================

alias aas='cd $HUB_SOURCES/sample-accelerators'

# vim: ft=zsh
