# etc/profile.d/90httpd
# =============================================================================
# HTTPD (Open Source builds of Apache httpd)
# =============================================================================

export HTTPD_QA_DIR=$ME_SRC_DIR/gitlab.eng.vmware.com/ccheetham/httpd-qa
alias -g pws=$HTTPD_QA_DIR

# vim: ft=zsh
