# etc/profile.d/90httpd
# =============================================================================
# Pivotal HTTPD (Open Source builds of Apache httpd)
# =============================================================================

alias pws='cd ~/src/appsuite/web-server-qa'
alias fhttpd='fly --target httpd'

# vim: ft=zsh
