# etc/profile.d/90me.ps1

$Env:ME_REPOS = "$Env:ME_SRC_DIR\github.com\ccheetham"

function meh { cd "$Env:ME_REPOS" }
function mes { cd "$Env:ME_REPOS\sandbox" }
