# etc/profile.d/90pivotal-steeltoe
# =============================================================================
# Pivotal Steeltoe
# =============================================================================

STEELTOE_REPOS=~/src/steeltoe

alias sts='cd $STEELTOE_REPOS/Samples'
alias stt='cd $STEELTOE_REPOS/Tooling'
alias std='cd $STEELTOE_REPOS/Dockerfiles'

# vim: ft=zsh
