# etc/profile.d/90pivotal-steeltoe
# =============================================================================
# Pivotal tc Server
# =============================================================================

STEELTOE_REPOS=~/src/repos/SteeltoeOSS

alias sts='cd $STEELTOE_REPOS/Samples'
alias stt='cd $STEELTOE_REPOS/Tooling'
alias std='cd $STEELTOE_REPOS/Dockerfiles'

# vim: ft=zsh
