# etc/profile.d/90pivotal-steeltoe
# =============================================================================
# Pivotal Steeltoe
# =============================================================================

STEELTOE_REPOS=~/src/steeltoe

alias std='cd $STEELTOE_REPOS/dockerfiles'
alias sti='cd $STEELTOE_REPOS/initializr'
alias sts='cd $STEELTOE_REPOS/samples'
alias stt='cd $STEELTOE_REPOS/tooling'

# vim: ft=zsh
