# etc/profile.d/90steeltoe
# =============================================================================
# Steeltoe
# =============================================================================

export STEELTOE_REPOS=~/src/github.com/SteeltoeOSS

alias st='cd $STEELTOE_REPOS'
alias stp='cd $STEELTOE_REPOS/Steeltoe'
alias sts='cd $STEELTOE_REPOS/Samples'
alias std='cd $STEELTOE_REPOS/Dockerfiles'
alias stt='cd $STEELTOE_REPOS/Tooling'
alias sti='cd $STEELTOE_REPOS/InitializrWorkspace'
alias stiw='cd $STEELTOE_REPOS/InitializrWorkspace/InitializrWeb'
alias stic='cd $STEELTOE_REPOS/InitializrWorkspace/InitializrConfig'
alias stiC='cd $STEELTOE_REPOS/InitializrWorkspace/InitializrConfigServer'
alias stik='cd $STEELTOE_REPOS/InitializrWorkspace/InitializrCli'
alias stis='cd $STEELTOE_REPOS/InitializrWorkspace/InitializrService'
alias ntt='cd $STEELTOE_REPOS/InitializrWorkspace/NetCoreToolTemplates'
alias nts='cd $STEELTOE_REPOS/InitializrWorkspace/NetCoreToolService'

# vim: ft=zsh
