# etc/profile.d/90steeltoe.ps1
# =============================================================================
# Steeltoe
# =============================================================================

$SteeltoeRepos="$Env:UserProfile\src\SteeltoeOSS"

function sts { cd "$SteeltoeRepos\Samples" }
function std { cd "$SteeltoeRepos\Dockerfiles" }
function stt { cd "$SteeltoeRepos\Tooling" }
function stil { cd "$SteeltoeRepos\Initializr" }
function stia { cd "$SteeltoeRepos\InitializrApi" }
function stiw { cd "$SteeltoeRepos\InitializrWeb" }
function stic { cd "$SteeltoeRepos\InitializrConfigServer" }
function stik { cd "$SteeltoeRepos\InitializrConfig" }
