# etc/profile.d/90steeltoe
# =============================================================================
# Steeltoe
# =============================================================================

export STEELTOE_REPOS=~/src/github.com/SteeltoeOSS
export STEELTOE_INCUBATOR_REPOS=~/src/github.com/steeltoeoss-incubator

alias -g st=$STEELTOE_REPOS
alias -g stl=$STEELTOE_REPOS/Steeltoe
alias -g sts=$STEELTOE_REPOS/Samples
alias -g std=$STEELTOE_REPOS/Dockerfiles
alias -g stt=$STEELTOE_REPOS/Tooling
alias -g stD=$STEELTOE_REPOS/Documentation
alias -g stMS=$STEELTOE_REPOS/MainSite
alias -g sti=$STEELTOE_REPOS/InitializrWorkspace
alias -g stiw=$STEELTOE_REPOS/InitializrWorkspace/InitializrWeb
alias -g stic=$STEELTOE_REPOS/InitializrWorkspace/InitializrConfig
alias -g stik=$STEELTOE_REPOS/InitializrWorkspace/InitializrConfigServer
alias -g stis=$STEELTOE_REPOS/InitializrWorkspace/InitializrService
alias -g ntt=$STEELTOE_REPOS/InitializrWorkspace/NetCoreToolTemplates
alias -g nts=$STEELTOE_REPOS/InitializrWorkspace/NetCoreToolService

alias -g stinc=$STEELTOE_INCUBATOR_REPOS
alias -g tt=$STEELTOE_INCUBATOR_REPOS/tanzu-thingies
alias -g dtg=$ME_REPOS/dotnet-dev-tap-gitops

ME_PATH+=:$STEELTOE_INCUBATOR_REPOS/tanzu-thingies/local/bin

# vim: ft=zsh
