# etc/profile.d/90steeltoe
# =============================================================================
# Steeltoe
# =============================================================================

STEELTOE_REPOS=~/src/SteeltoeOSS

alias std='cd $STEELTOE_REPOS/Dockerfiles'
alias stil='cd $STEELTOE_REPOS/Initializr'
alias stia='cd $STEELTOE_REPOS/InitializrApi'
alias stiw='cd $STEELTOE_REPOS/InitializrWeb'
alias sts='cd $STEELTOE_REPOS/Samples'
alias stt='cd $STEELTOE_REPOS/Tooling'

# vim: ft=zsh
