# etc/profile.d/90steeltoe
# =============================================================================
# Steeltoe
# =============================================================================

alias stil='Initializr'
alias stia='InitializrApi'
alias stiw='InitializrWeb'
alias stic='Steeltoe.Initializr.ConfigServer'
alias stiC='Steeltoe.Initializr.Configuration'
alias sts='Samples'
alias stt='Tooling'

# vim: ft=zsh
