# etc/profile.d/90steeltoe
# =============================================================================
# Steeltoe
# =============================================================================

STEELTOE_REPOS=~/src/SteeltoeOSS

alias std='cd $STEELTOE_REPOS/Dockerfiles'
alias sti='cd $STEELTOE_REPOS/initializr'
alias sts='cd $STEELTOE_REPOS/Samples'
alias stt='cd $STEELTOE_REPOS/Tooling'

# vim: ft=zsh
