# etc/profile.d/90tanzu-starter-service
# =============================================================================
# Tanzu Starter Service
# =============================================================================

export TSS_SRC=~/src/gitlab.eng.vmware.com/tanzu-starter-service/tss

alias tss="java -jar $TSS_SRC/tss-cli/target/tss-cli-1.0.0-SNAPSHOT.jar"
alias tss-server="java -jar $TSS_SRC/tss-server/target/tss-server-1.0.0-SNAPSHOT.jar"

# vim: ft=zsh
