# etc/profile.d/90tanzu
# =============================================================================
# Tanzu
# -----------------------------------------------------------------------------
# Tanzu Application Platform
# https://tanzu.vmware.com/developer/guides/tanzu-application-platform-local-devloper-install/
# =============================================================================

TANZU_CLI_VERSION=${TANZU_CLI_VERSION:-v0.11.2}

# vim: ft=zsh
