# etc/profile.d/90tanzu
# =============================================================================
# Tanzu
# -----------------------------------------------------------------------------
# Tanzu Application Platform
# https://tanzu.vmware.com/developer/guides/tanzu-application-platform-local-developer-install/
# =============================================================================

TANZU_CLI_VERSION=${TANZU_CLI_VERSION:-v0.11.2}
TANZU_REPOS=$ME_SRC_DIR/github.com/vmware-tanzu

alias vt='cd $TANZU_REPOS'
alias vtvs='cd $TANZU_REPOS/tanzu-developer-tools-for-visual-studio'
alias vtc='cd $TANZU_REPOS/tanzu-developer-tools-for-vscode'
alias vtij='cd $TANZU_REPOS/tanzu-developer-tools-for-intellij'
alias vtls='cd $TANZU_REPOS/tanzu-developer-tools-language-server'

# vim: ft=zsh
