# etc/profile.d/90tanzu
# =============================================================================
# Tanzu
# -----------------------------------------------------------------------------
# Tanzu Application Platform
# https://tanzu.vmware.com/developer/guides/tanzu-application-platform-local-developer-install/
# =============================================================================

TANZU_REPOS=$ME_SRC_DIR/github.com/vmware-tanzu

alias -g vt=$TANZU_REPOS
alias -g dtc=$TANZU_REPOS/tanzu-developer-tools-for-vscode
alias -g dtj=$TANZU_REPOS/tanzu-developer-tools-for-intellij
alias -g dtl=$TANZU_REPOS/tanzu-developer-tools-language-server
alias -g dts=$TANZU_REPOS/tanzu-developer-tools-for-visual-studio
alias -g aas=$TANZU_REPOS/application-accelerator-samples

# vim: ft=zsh
