# etc/profile.d/90tcserver.ps1
# =============================================================================
# tc Server
# =============================================================================

$TcsQa="$Env:UserProfile\src\github.com\appsuite\tcs-tc-server-qa"

function tcd {
  cd "$TcsQa"
}
