# etc/profile.d/92cloudfoundry
# =============================================================================
# Cloud Foundry
# =============================================================================

CF_REPOS=$ME_SRC_DIR/github.com/pivotal-cf

alias -g cfr=$CF_REPOS
alias -g kbp=$CF_REPOS/kerberos-buildpack

# vim: ft=zsh
