# etc/profile.d/92tanzu
# =============================================================================
# VMware Tanzu
# -----------------------------------------------------------------------------
# VMware Tanzu Application Platform
# https://tanzu.vmware.com/developer/guides/tanzu-application-platform-local-developer-install/
# =============================================================================

TANZU_REPOS=$ME_SRC_DIR/github.com/vmware-tanzu
ME_PATH+=:$STEELTOE_INCUBATOR_REPOS/tanzu-thingies/bin

alias -g vtr=$TANZU_REPOS
alias -g dtc=$TANZU_REPOS/tanzu-developer-tools-for-vscode
alias -g dtj=$TANZU_REPOS/tanzu-developer-tools-for-intellij
alias -g dtl=$TANZU_REPOS/tanzu-developer-tools-language-server
alias -g dts=$TANZU_REPOS/tanzu-developer-tools-for-visual-studio
alias -g aas=$TANZU_REPOS/application-accelerator-samples

alias tzl="tz apps workload list"
alias tzg="tz apps workload get"
alias tzd="tz apps workload delete"
alias tzt="tz apps workload tail --timestamp"


tz() {
  echo % tanzu $*
  tanzu $*
}
# vim: ft=zsh
