# etc/profile.d/95mepath.ps1
# =============================================================================
# finalize path setups
# =============================================================================

$env:PATH = "$ME_PATH;$env:PATH"
